# ipset-fw.timer
[Unit]
Description=update ipsets timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
